# vskeddemos/configbackup/host/nacos/cluster/nacos1/nacos1.service
[Unit]
Description=nacos
After=network.target
 
[Service]
Type=forking
ExecStart=/opt/nacos/nacos1/bin/startup.sh
ExecReload=/opt/nacos/nacos1/bin/shutdown.sh
ExecStop=/opt/nacos/nacos1/bin/shutdown.sh
PrivateTmp=true
 
[Install]  
WantedBy=multi-user.target
